# panel/systemd/tenodera-gateway.service
[Unit]
Description=Tenodera Web Console Gateway
Documentation=https://github.com/tenodera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tenodera-gateway
Environment=RUST_LOG=tenodera_gateway=info
Environment=TENODERA_BRIDGE_BIN=/usr/local/bin/tenodera-bridge
Environment=TENODERA_UI_DIR=/usr/share/tenodera/ui
# TLS — uncomment and set paths for production:
# Environment=TENODERA_TLS_CERT=/etc/tenodera/tls/cert.pem
# Environment=TENODERA_TLS_KEY=/etc/tenodera/tls/key.pem
# Environment=TENODERA_ALLOW_UNENCRYPTED=0

# Security hardening
ProtectSystem=strict
# Bridge child process needs write access for user/group management (useradd, etc.)
# and networking configuration. /etc is needed for passwd/shadow/group lock files.
ReadWritePaths=/etc /var/log /home /var/mail
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
# useradd/groupadd create lock files with specific permissions
RestrictSUIDSGID=no

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
