# panel/systemd/tenodera-gateway.service
[Unit]
Description=Tenodera Web Console Gateway
Documentation=https://github.com/tenodera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tenodera-gateway
EnvironmentFile=-/etc/tenodera/gateway.env

# Security hardening
# NOTE: ProtectSystem=strict is NOT used here because tenodera-bridge
# (spawned as a child process) must write to /usr (package installs),
# /var (dpkg/apt/dnf state, docker, systemd), /run (systemd, docker
# sockets), and arbitrary paths via PTY sessions. Enumerating all
# required ReadWritePaths would effectively whitelist the entire
# filesystem, defeating the purpose of ProtectSystem.
# Other hardening directives below still provide meaningful protection.
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
# useradd/groupadd create lock files with specific permissions
RestrictSUIDSGID=no

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
