# Tenodera Agent/systemd/tenodera-agent.service
[Unit]
Description=Tenodera Agent — system management daemon
Documentation=https://github.com/tenodera/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tenodera-agent
Restart=on-failure
RestartSec=5

# Config file location
Environment=TENODERA_AGENT_CONFIG=/etc/tenodera/agent.toml

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tenodera-agent

[Install]
WantedBy=multi-user.target
